# vm-systemd/selinux-autolabel-mark.service.d/30_qubes.conf
[Unit]
After=rw.mount

[Service]
ExecStart=-/bin/rm -f /rw/.autorelabel
