# vm-systemd/systemd-timesyncd.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/clocksync
After=qubes-sysinit.service
