# vm-systemd/tor.service.d/30_qubes.conf
[Unit]
ConditionPathExists=!/var/run/qubes/this-is-templatevm
