# vm-systemd/chronyd.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/ntpd
