# vm-systemd/ntpd.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/ntpd
After=qubes-sysinit.service
