# vm-systemd/polkit.service.d/30_qubes.conf
[Unit]
ConditionPathExists=!/var/run/qubes-service/minimal-usbvm
