# vm-systemd/cron.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/crond

[Service]
After=qubes-sysinit.service
