# vm-systemd/avahi-daemon.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/avahi
