# vm-systemd/exim4.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/exim4
