# vm-systemd/qubes-rootfs-resize.service.d/30_qubes-kvm.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/hypervisor-kvm
# Device is 'dev-vda' in KVM (dev-xvda in  XEN)
After=
After=qubes-sysinit.service dev-vda.device
