# vm-systemd/cups.socket
.include /lib/systemd/system/cups.socket
[Unit]
ConditionPathExists=/var/run/qubes-service/cups
