# vm-systemd/haveged.service
.include /lib/systemd/system/haveged.service

[Install]
WantedBy=
WantedBy=multi-user.target
