# vm-systemd/home.mount
[Mount]
What=/rw/home
Where=/home
Type=none
Options=noauto,bind,defaults,nosuid,nodev
