# vm-systemd/qubes-sync-time.timer
[Unit]
Description=Update system time each 6h
ConditionPathExists=!/var/run/qubes-service/clocksync

[Timer]
OnBootSec=10s
OnUnitActiveSec=6h


