# vm-systemd/usr-local.mount
[Mount]
What=/rw/usrlocal
Where=/usr/local
Type=none
Options=noauto,bind,defaults
