# vm-systemd/user/gvfs-daemon.service.d/30_qubes.conf
[Unit]
ConditionPathExists=!/var/run/qubes-service/minimal-netvm
ConditionPathExists=!/var/run/qubes-service/minimal-usbvm
