# debian/qubes-core-agent-passwordless-root.postinst
#!/bin/bash
# postinst script for qubes-core-agent-passwordless-root
#
# see: dh_installdeb(1)

set -e

# The postinst script may be called in the following ways:
#   * <postinst> 'configure' <most-recently-configured-version>
#   * <old-postinst> 'abort-upgrade' <new version>
#   * <conflictor's-postinst> 'abort-remove' 'in-favour' <package>
#     <new-version>
#   * <postinst> 'abort-remove'
#   * <deconfigured's-postinst> 'abort-deconfigure' 'in-favour'
#     <failed-install-package> <version> 'removing'
#     <conflicting-package> <version>
#
#    For details, see http://www.debian.org/doc/debian-policy/ or
# https://www.debian.org/doc/debian-policy/ch-maintainerscripts.html or
# the debian-policy package

pam-auth-update --package

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

# If systemd is running try starting qubes-admin-authzd even if the user
# blocked the normal start by dh (above) via policy-rc.d. Not running it will
# break passwordless auth, so this is most likely not what the user wanted. But
# if they really want to they still can mask the service.

if [[ -e /run/systemd/system ]] &&
   ! systemctl is-active --quiet qubes-admin-authzd.service &&
   systemctl is-enabled --quiet qubes-admin-authzd.service
then
    systemctl start qubes-admin-authzd.service || true
fi

exit 0

# vim: set ts=4 sw=4 sts=4 et :
