# vm-systemd/cups.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/cups
After=qubes-sysinit.service

[Service]
ExecStartPost=/usr/lib/qubes/init/control-printer-icon.sh
