# vm-systemd/user/tracker-miner-fs-3.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/run/qubes-service/tracker
