# vm-systemd/qubes-sysinit.service.d/30_qubes-kvm.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/hypervisor-kvm
# No depend on 'proc-xen.mount' in KVM
After=
After=systemd-modules-load.service qubes-db.service
