# vm-systemd/user/evolution-source-registry.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/run/qubes-service/evolution-data-server
