# vm-systemd/rescue.service.d/30_qubes.conf
[Service]
# Ensure that the console is the secure Xen console,
# not e.g. a serial console that is exposed to the outside world.
TTYPath=/dev/hvc0
Environment=SYSTEMD_SULOGIN_FORCE=1
