# vm-systemd/systemd-random-seed.service.d/30_qubes.conf
[Unit]
After=qubes-db.service
Before=sysinit.target

[Service]
ExecStartPre=/usr/lib/qubes/init/qubes-random-seed.sh
