# vm-systemd/firewalld.service.d/30_qubes.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/firewalld
After=qubes-firewall.service
