# vm-systemd/serial-getty@.service.d/30_qubes.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --login-pause --keep-baud 115200,38400,9600 %I $TERM
