# vm-systemd/ModemManager.service.d/30_qubes.conf
[Unit]
ConditionPathExists=|/var/run/qubes-service/network-manager
ConditionPathExists=|/var/run/qubes-service/modem-manager
After=qubes-sysinit.service
